# scripts/mystartup.service
[Unit]
Description=Custom Pop OS startup script

[Service]
Type=simple
ExecStart=/bin/bash /usr/bin/startup.sh

[Install]
WantedBy=multi-user.target
